# crontab
TZ=Asia/Seoul
0 0 * * * /usr/local/bin/python /app/services/daily_tasks.py >> /app/cron.log 2>&1
